# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:shellinaboxd]
command=/usr/sbin/shellinaboxd -u root -b
